# templates/default/sinopia.service.erb
# Automatically Generated by Chef

[Unit]
Description=Sinopia NPM cache & private repo server

[Service]
ExecStart=/usr/bin/sinopia -c <%= File.join(node['sinopia']['confdir'], 'config.yaml') %>
Type=simple
User=<%= node['sinopia']['user'] %>
Group=<%= node['sinopia']['group'] %>

[Install]
WantedBy=multi-user.target
